# .docker/php/crontab
* * * * * /usr/local/bin/run-cron.sh >> /var/www/storage/logs/cron.log 2>&1
